# bin/gemfire-locator
#!/bin/bash
###
# chkconfig: 235 98 55
# description: Manages the Gemfire locator service
###

OLDPWD=$PWD
cd `dirname $0`
BIN=`pwd`;BASE=`dirname $BIN`
cd $OLDPWD

source $BASE/bin/common.sh
source $BASE/bin/setenv.sh

case "$1" in
  start)
        echo "Starting $TAG gemfire locator"
		export CLASSPATH
        $GEMFIRE/bin/gemfire start-locator \
			-dir=$GEMFIRE_WORK -port=$LOCATOR_PORT \
			-DgemfirePropertyFile=$LOCATOR_PROP_FILE \
			-Dgemfire.log-file=$LOCATOR_LOG \
			-Denv=$TAG
        ;;
  stop)
        echo "Stopping $TAG gemfire locators"
        $GEMFIRE/bin/gemfire stop-locator \
			-dir=$GEMFIRE_WORK -port=$LOCATOR_PORT
        ;;

  *)
        echo "Usage: /sbin/service gemfire-locator {start|stop}"
        exit 1
esac

exit 0
